# universal_controller/rootfs/etc/services.d/universal-controller/run
#!/usr/bin/with-contenv bashio

# Get options from configuration
declare log_level
declare port

log_level=$(bashio::config 'log_level')
port=$(bashio::config 'port')

bashio::log.info "Starting Universal Controller..."
bashio::log.info "Log level: ${log_level}"
bashio::log.info "Port: ${port}"

# Set environment variables
export LOG_LEVEL="${log_level}"
export PORT="${port}"

# Change to app directory
cd /app

# Start the Node.js server
exec node server.js
